# conf.d/telegram.conf
[program:telegram]
command=telegram
autorestart=true
